# threat_hunts/macos/persistence_mechanisms.yml
---
checks:
  - id: "TH-MAC-001"
    title: "Detect suspicious launch agents and daemons"
    description: "Identifies potentially malicious launch agents and daemons that could be used for persistence"
    severity: "high"
    query: "SELECT path, name, program, program_arguments FROM launchd WHERE path LIKE '/Library/LaunchAgents/%' OR path LIKE '/Library/LaunchDaemons/%' OR path LIKE '/Users/%/Library/LaunchAgents/%';"
    condition:
      type: "manual_review"
      field: "all"
      value: "N/A"

  - id: "TH-MAC-002"
    title: "Detect suspicious login items"
    description: "Identifies potentially malicious login items that could be used for persistence"
    severity: "high"
    query: "SELECT username, name, path FROM users JOIN login_items USING (uid);"
    condition:
      type: "manual_review"
      field: "all"
      value: "N/A"

  - id: "TH-MAC-003"
    title: "Detect suspicious cron jobs"
    description: "Identifies potentially malicious cron jobs that could be used for persistence"
    severity: "high"
    query: "SELECT c.command, c.path, u.username FROM crontab c JOIN users u ON c.uid = u.uid WHERE c.command LIKE '%nc %' OR c.command LIKE '%netcat%' OR c.command LIKE '%base64 -d%' OR c.command LIKE '%curl%|bash%' OR c.command LIKE '%wget%|bash%';"
    condition:
      type: "equals"
      field: "count"
      value: "0"

  - id: "TH-MAC-004"
    title: "Detect suspicious kernel extensions"
    description: "Identifies potentially malicious kernel extensions that could be used for rootkits"
    severity: "critical"
    query: "SELECT name, version, linked_against, path, loadable FROM kernel_extensions WHERE NOT path LIKE '/System/Library/Extensions/%' AND NOT path LIKE '/Library/Apple/%';"
    condition:
      type: "manual_review"
      field: "all"
      value: "N/A"

  - id: "TH-MAC-005"
    title: "Detect suspicious startup items"
    description: "Identifies potentially malicious startup items that could be used for persistence"
    severity: "high"
    query: "SELECT name, path, type, source FROM startup_items WHERE path NOT LIKE '/System/%' AND path NOT LIKE '/Library/Apple/%';"
    condition:
      type: "manual_review"
      field: "all"
      value: "N/A"

  - id: "TH-MAC-006"
    title: "Detect suspicious processes"
    description: "Identifies potentially malicious processes with suspicious command lines"
    severity: "high"
    query: "SELECT p.pid, p.name, p.path, p.cmdline, p.uid, u.username FROM processes p JOIN users u ON p.uid = u.uid WHERE p.cmdline LIKE '%nc %' OR p.cmdline LIKE '%netcat%' OR p.cmdline LIKE '%base64 -d%' OR p.cmdline LIKE '%curl%|bash%' OR p.cmdline LIKE '%wget%|bash%' OR p.cmdline LIKE '%python%socket%' OR p.cmdline LIKE '%perl%socket%';"
    condition:
      type: "equals"
      field: "count"
      value: "0"

  - id: "TH-MAC-007"
    title: "Detect suspicious network connections"
    description: "Identifies potentially malicious network connections to unusual ports or addresses"
    severity: "high"
    query: "SELECT s.pid, p.name, p.path, p.cmdline, s.local_address, s.local_port, s.remote_address, s.remote_port, s.state FROM process_open_sockets s JOIN processes p ON s.pid = p.pid WHERE (s.remote_port = 4444 OR s.remote_port = 4545 OR s.remote_port = 1337 OR s.remote_port = 31337 OR s.remote_port = 666 OR s.remote_port = 6666 OR s.remote_port = 6667 OR s.remote_port = 6668 OR s.remote_port = 6669 OR s.remote_port = 9999) AND s.state = 'ESTABLISHED';"
    condition:
      type: "equals"
      field: "count"
      value: "0"

  - id: "TH-MAC-008"
    title: "Detect suspicious dyld hijacking"
    description: "Identifies potential dyld hijacking in system directories"
    severity: "high"
    query: "SELECT path, directory, filename FROM file WHERE path LIKE '/usr/local/lib/%' OR path LIKE '/usr/lib/%' AND filename LIKE '%.dylib';"
    condition:
      type: "manual_review"
      field: "all"
      value: "N/A"

  - id: "TH-MAC-009"
    title: "Detect suspicious browser extensions"
    description: "Identifies potentially malicious browser extensions"
    severity: "medium"
    query: "SELECT * FROM safari_extensions UNION SELECT * FROM chrome_extensions;"
    condition:
      type: "manual_review"
      field: "all"
      value: "N/A"

  - id: "TH-MAC-010"
    title: "Detect suspicious system modifications"
    description: "Identifies modifications to system integrity protection and other security features"
    severity: "critical"
    query: "SELECT * FROM sip_config WHERE enabled = 0;"
    condition:
      type: "equals"
      field: "count"
      value: "0"
